# cutie-ui-io.service.d/10-desktop.conf
[Service]
Environment=QT_IM_MODULE=ibus
Environment=CUTIE_DESKTOP=1
Environment=CUTIE_DEFAULT_MODE=desktop