# cutie-ui-io.service.d/20-convergence.conf
[Service]
Environment=CUTIE_DESKTOP=1
Environment=CUTIE_MOBILE=1