# cutie-ui-io.service.d/30-hybris.conf
[Unit]
After=lxc@android.service

[Service]
Environment=QT_QPA_PLATFORM=hwcomposer
Environment=EGL_PLATFORM=hwcomposer
Environment=QT_SCALE_FACTOR="3"
ExecStartPre=/usr/bin/waitforservice init.svc.vendor.hwcomposer-*
User=32011
