# android-service@hwcomposer.service.d/20-cutie.conf
[Unit]
BindsTo=cutie-ui-io.service
Before=cutie-ui-io.service
